# .config/nvim/minio/server-configurations/arch-desktop/crontab
# run every 20 minutes to restart scdaemon
*/20 * * * * kev /home/kev/.local/scripts/restart-scdaemon
